# interface-api.service
[Unit]
Description=Interface API Service
After=network.target
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
User=interface-api
Group=interface-api
WorkingDirectory=/opt/interface-api
ExecStart=/opt/interface-api/bin/api
Restart=on-failure
RestartSec=5

EnvironmentFile=/opt/interface-api/.env.default
EnvironmentFile=/opt/interface-api/.env
Environment="XDG_CACHE_HOME=/opt/interface-api/.cache"
Environment="GOCACHE=/opt/interface-api/.cache/go-build"
Environment="GOMODCACHE=/opt/interface-api/.cache/go-mod"

; Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/interface-api
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
